# build/root/etc/s6-overlay/s6-rc.d/bootloader/run
#!/command/with-contenv bash

source /usr/bin/debug.sh

echo -e ""
figlet -f ansi_shadow.flf "SWAT4"
echo "A Docker container by MisterCalvin (v$(cat "/container/.CONTAINER_VERSION")) | discord.sbotnas.io"