# build/root/etc/s6-overlay/s6-rc.d/swat4-svc/run
#!/bin/sh
exec simple_logger -l "/container/swat4/swat4-container.log" "swat4-server" "true" s6-setuidgid wine /usr/bin/swat4_startup.sh 2>&1
