# build/root/etc/s6-overlay/s6-rc.d/swat4-svc/run
#!/command/execlineb -P
s6-setuidgid wine 
"/usr/bin/swat4.sh"
