# build/root/etc/s6-overlay/s6-rc.d/init-user/run
#!/command/with-contenv bash

PUID=${PUID:-1000}
PGID=${PGID:-1000}

if [ ! "$PUID" = "1000" ] || [ ! "$PGID" = "1000" ]; then
    	echo "New UID/GID specified, modifying current container user"
	echo "Please note - this does not change file permissions on the mounted volume! This will only chown /container and /container/.wine to the specified UID/GID"
	echo "If the server is crashing or refusing to boot, check file permissions on the host mount!"
	groupmod --gid "$PGID" wine
    	usermod --uid "$PUID" --gid "$PGID" wine
	chown "$PUID:$PGID" /container 
fi

if [ -z "$(ls -A /container/.wine)" ]; then
	printf "true" > /var/run/s6/container_environment/FIRST_RUN
else
	printf "false" > /var/run/s6/container_environment/FIRST_RUN
fi

chown -R $PUID:$PGID /container/.wine